# deploy/systemd/voice-cloning-api.service
[Unit]
Description=Voice Cloning API Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/voice-cloning-using-qwen3
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/ubuntu/voice-cloning-using-qwen3/scripts/run_api.py --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
