# deploy/fastapi.service
[Unit]
Description=FastAPI Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/actions-runner/_work/voice-cloning-using-qwen3/voice-cloning-using-qwen3

# Point directly at the venv's Python — not /usr/bin/python3
ExecStart=/home/ubuntu/actions-runner/_work/voice-cloning-using-qwen3/voice-cloning-using-qwen3/venv/bin/python scripts/run_api.py --host 127.0.0.1 --port 8000

Environment=PYTHONUNBUFFERED=1

Restart=always
RestartSec=5
# Add a restart limit so it stops hammering if it keeps failing
StartLimitIntervalSec=60
StartLimitBurst=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target