# deploy/fastapi.service
[Unit]
Description=FastAPI Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/actions-runner/_work/voice-cloning-using-qwen3/voice-cloning-using-qwen3

# Use venv Python directly — no need to activate
ExecStart=/home/ubuntu/actions-runner/_work/voice-cloning-using-qwen3/voice-cloning-using-qwen3/venv/bin/python scripts/run_api.py --host 127.0.0.1 --port 8000

# Explicitly set PATH and PYTHONPATH for safety
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/ubuntu/actions-runner/_work/voice-cloning-using-qwen3/voice-cloning-using-qwen3/venv/bin:/usr/bin:/bin

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target